# MagicMemory/extern/etc/udev/rules.d/magicmemory.rules
KERNEL=="loop0", ACTION=="add", RUN+="/magic_udev"

