# MiniClock/etc/udev/rules.d/miniclock.rules
KERNEL=="loop0", ACTION=="add", RUN+="/usr/local/MiniClock/miniclock.sh"
