# ScreenSaver/etc/udev/rules.d/screensaver.rules
KERNEL=="loop0", ACTION=="add", RUN+="/usr/local/ScreenSaver/screensaver.sh"
