# ScreenSaver/etc/udev/rules.d/screensaver.rules
KERNEL=="loop0", RUN+="/usr/local/ScreenSaver/screensaver.sh"
