# AutoPatch/etc/udev/rules.d/autopatch.rules
KERNEL=="fsl-usb2-udc", ACTION=="remove", RUN+="/usr/local/AutoPatch/autopatch.sh"
