# AutoShelf/etc/udev/rules.d/autoshelf.rules
KERNEL=="fsl-usb2-udc", ACTION=="remove", RUN+="/usr/local/AutoShelf/autoshelf.sh"
