# AutoShelf/etc/udev/rules.d/autoshelf.rules
KERNEL=="fsl-usb2-udc", SUBSYSTEM=="drivers", RUN+="/usr/local/AutoShelf/autoshelf.sh"
