# NickelIcons/etc/udev/rules.d/nickelicons.rules
KERNEL=="fsl-usb2-udc", ACTION=="remove", RUN+="/usr/local/NickelIcons/nickelicons.sh"
