# NickelIcons/etc/udev/rules.d/nickelicons.rules
SUBSYSTEM=="module", DEVPATH=="/module/g_file_storage", ACTION=="remove", RUN+="/usr/local/NickelIcons/nickelicons.sh"
