# fatpatch/etc/udev/rules.d/fatpatch.rules
KERNEL=="loop0", RUN+="/fatpatch/fatpatch.sh"
