# fatpatch/etc/udev/rules.d/fatpatch.rules
KERNEL=="loop0", ACTION=="add", RUN+="/fatpatch/fatpatch.sh"
