# WebPortal/etc/udev/rules.d/webportal.rules
KERNEL=="loop0", RUN+="/usr/local/WebPortal/webportal.sh"
