# WebPortal/etc/udev/rules.d/webportal.rules
KERNEL=="loop0", ACTION=="add", RUN+="/usr/local/WebPortal/webportal.sh"
