# UsbDebug/etc/udev/rules.d/usbdebug.rules
KERNEL=="fsl-usb2-udc", SUBSYSTEM=="drivers", RUN+="/usr/local/UsbDebug/usbdebug.sh"
