# UsbDebug/etc/udev/rules.d/usbdebug.rules
SUBSYSTEM=="module", DEVPATH=="/module/g_file_storage", RUN+="/usr/local/UsbDebug/usbdebug.sh"
