# FlashCard/etc/udev/rules.d/flashcard.rules
KERNEL=="loop0", ACTION=="add", RUN+="/usr/local/FlashCard/flashcard.sh"
