# TelnetFtp/etc/udev/rules.d/telnetftp.rules
KERNEL=="loop0", ACTION=="add", RUN+="/sbin/my-telnetftp"

