# disablesickel/etc/udev/rules.d/disablesickel.rules
KERNEL=="loop0", RUN+="/disablesickel.sh"
