# etc/udev/rules.d/autopatch.rules
# udev runs the startup script for us
KERNEL=="loop0", RUN+="/mod/AutoPatch/autopatch"
